# TF-algorithm/InputData/OnPremData/salt/analysis20231010094514/code_health_details/f1b195fed82feb9ebacea7ac871ba01d7d883716.json
{"last-year-high-resolution-score":null,"high-resolution-score":6.674895159827102,"revision":"b24911a3e1","calculated":true,"name":"salt/salt/modules/osquery.py","last-month":6,"last-month-details":{"name":"salt/salt/modules/osquery.py","score":6,"high-resolution-score":6.674895159827102,"unbiased-score":6.674895159827102,"details":{"cc-mean":1.6666666666666667,"main-body-cc":0,"lines-in-file":715,"complex-functions":[{"name":"_osquery_cmd","start-line":65,"end-line":106,"cc":12}],"cc-median":1.0,"nested":{"contains-named-nested-functions":true,"nested-complexity-of-interest":[{"name":"_osquery_cmd","start-line":65,"end-line":106,"nested-complexity":5}],"max-nested-in-global-scope":0,"max-nested-complexity-depth":5,"max-nested-complexity-depth-name":"_osquery_cmd"},"excess-long-functions":[],"fns-with-duplicated-code":[{"name":"rpm_packages","start-line":124,"end-line":136},{"name":"kernel_integrity","start-line":139,"end-line":151},{"name":"kernel_modules","start-line":154,"end-line":166},{"name":"memory_map","start-line":169,"end-line":181},{"name":"process_memory_map","start-line":184,"end-line":196},{"name":"shared_memory","start-line":199,"end-line":211},{"name":"apt_sources","start-line":214,"end-line":226},{"name":"deb_packages","start-line":229,"end-line":241},{"name":"alf","start-line":647,"end-line":659},{"name":"alf_exceptions","start-line":662,"end-line":674},{"name":"alf_services","start-line":692,"end-line":704},{"name":"apps","start-line":707,"end-line":719},{"name":"certificates","start-line":722,"end-line":734},{"name":"homebrew_packages","start-line":767,"end-line":779},{"name":"launchd","start-line":842,"end-line":854},{"name":"nfs_shares","start-line":857,"end-line":869},{"name":"nvram","start-line":872,"end-line":884},{"name":"preferences","start-line":887,"end-line":899},{"name":"quarantine","start-line":902,"end-line":914},{"name":"startup_items","start-line":932,"end-line":944},{"name":"kernel_extensions","start-line":812,"end-line":824},{"name":"firefox_addons","start-line":752,"end-line":764},{"name":"iokit_devicetree","start-line":782,"end-line":794},{"name":"iokit_registry","start-line":797,"end-line":809},{"name":"keychain_items","start-line":827,"end-line":839},{"name":"xattr_where_from","start-line":947,"end-line":959},{"name":"xprotect_entries","start-line":962,"end-line":974},{"name":"xprotect_reports","start-line":977,"end-line":989},{"name":"alf_explicit_auths","start-line":677,"end-line":689},{"name":"chrome_extensions","start-line":737,"end-line":749},{"name":"safari_extensions","start-line":917,"end-line":929}],"complex-conditionals":[{"name":"version","details":[{"start-line":120,"end-line":121,"branches":2}]}],"longest-fn-loc-name":"_osquery_cmd","n-functions":75,"active-code-size":701,"bumps":{},"longest-fn-loc":37,"cohesion":1,"clone-ratio":0.4133333333333333,"fn-args":{"ctors-with-many-args":[],"max-args":4,"n-primitives":0,"fns-with-many-args":[],"max-ctor-args":0,"n-string-args":0,"mean-args":1.9333333333333333,"max-ctor-arg-name":"<unknown>","n-args":145,"max-arg-name":"_osquery_cmd"},"cc-max":12,"n-clones":31,"congestion":{"authors":1,"fractal-value":0.0},"median-fn-loc":8.0,"cc-max-name":"_osquery_cmd","cc-total":125},"revision":"b24911a3e1","date":"2015-07-08"},"last-year-details":{"details":{"lines-in-file":null},"revision":null,"date":null},"unbiased-score":6.674895159827102,"details":{"cc-mean":1.6666666666666667,"main-body-cc":0,"lines-in-file":715,"social":{"owner":"Gareth J. Greenaway","ownership":1.0,"knowledge-loss":0.0},"complex-functions":[{"name":"_osquery_cmd","start-line":65,"end-line":106,"cc":12}],"cc-median":1.0,"nested":{"contains-named-nested-functions":true,"nested-complexity-of-interest":[{"name":"_osquery_cmd","start-line":65,"end-line":106,"nested-complexity":5}],"max-nested-in-global-scope":0,"max-nested-complexity-depth":5,"max-nested-complexity-depth-name":"_osquery_cmd"},"code-comment-match":[{"n-matches":0,"pattern-name":"Detect TODOs"}],"excess-long-functions":[],"fns-with-duplicated-code":[{"name":"rpm_packages","start-line":124,"end-line":136},{"name":"kernel_integrity","start-line":139,"end-line":151},{"name":"kernel_modules","start-line":154,"end-line":166},{"name":"memory_map","start-line":169,"end-line":181},{"name":"process_memory_map","start-line":184,"end-line":196},{"name":"shared_memory","start-line":199,"end-line":211},{"name":"apt_sources","start-line":214,"end-line":226},{"name":"deb_packages","start-line":229,"end-line":241},{"name":"alf","start-line":647,"end-line":659},{"name":"alf_exceptions","start-line":662,"end-line":674},{"name":"alf_services","start-line":692,"end-line":704},{"name":"apps","start-line":707,"end-line":719},{"name":"certificates","start-line":722,"end-line":734},{"name":"homebrew_packages","start-line":767,"end-line":779},{"name":"launchd","start-line":842,"end-line":854},{"name":"nfs_shares","start-line":857,"end-line":869},{"name":"nvram","start-line":872,"end-line":884},{"name":"preferences","start-line":887,"end-line":899},{"name":"quarantine","start-line":902,"end-line":914},{"name":"startup_items","start-line":932,"end-line":944},{"name":"kernel_extensions","start-line":812,"end-line":824},{"name":"firefox_addons","start-line":752,"end-line":764},{"name":"iokit_devicetree","start-line":782,"end-line":794},{"name":"iokit_registry","start-line":797,"end-line":809},{"name":"keychain_items","start-line":827,"end-line":839},{"name":"xattr_where_from","start-line":947,"end-line":959},{"name":"xprotect_entries","start-line":962,"end-line":974},{"name":"xprotect_reports","start-line":977,"end-line":989},{"name":"alf_explicit_auths","start-line":677,"end-line":689},{"name":"chrome_extensions","start-line":737,"end-line":749},{"name":"safari_extensions","start-line":917,"end-line":929}],"complex-conditionals":[{"name":"version","details":[{"start-line":120,"end-line":121,"branches":2}]}],"longest-fn-loc-name":"_osquery_cmd","n-functions":75,"active-code-size":701,"bumps":{},"longest-fn-loc":37,"cohesion":1,"delta":{"active-code-delta":0,"cc-total-delta":0,"nested-depth-delta":0,"n-functions-delta":0},"clone-ratio":0.4133333333333333,"fn-args":{"ctors-with-many-args":[],"max-args":4,"n-primitives":0,"fns-with-many-args":[],"max-ctor-args":0,"n-string-args":0,"mean-args":1.9333333333333333,"max-ctor-arg-name":"<unknown>","n-args":145,"max-arg-name":"_osquery_cmd"},"cc-max":12,"n-clones":31,"congestion":{"authors":0,"fractal-value":0.0},"median-fn-loc":8.0,"cc-max-name":"_osquery_cmd","cc-total":125},"score":6,"last-month-high-resolution-score":6.674895159827102,"file-details":{"name":"salt/salt/modules/osquery.py","resolved":"/Users/andreaskarlsson/Documents/Exjobb/algorithm/test-oracle/repos/salt/salt/modules/osquery.py","repo-path":"/Users/andreaskarlsson/Documents/Exjobb/algorithm/test-oracle/repos/salt","relative-name":"salt/modules/osquery.py"},"last-year":0}